# SOURCES/unisonsync.service
[Unit]
Description=Sync user's home with unison

[Service]
Type=oneshot
ExecStart=/usr/bin/sh /usr/lib/roaming-profiles/unisonsync.sh

